# roles/install_springapp/templates/springapp.service.j2
[Unit]
Description=webserver Daemon

[Service]
User=root
ExecStart={{ exec_start }}

[Install]
WantedBy=multi-user.target